# scripts/com.user.backup-mailru.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.user.backup-mailru</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/evgenijgrudev/scripts/backup_mailru.sh</string>
    </array>

    <!-- Запуск каждый день в 11:00 -->
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>11</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <!-- Если Mac спал в 11:00 — запустит при следующем пробуждении -->

    <key>StandardOutPath</key>
    <string>/Users/evgenijgrudev/.backup_mailru_logs/launchd_stdout.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/evgenijgrudev/.backup_mailru_logs/launchd_stderr.log</string>

    <key>RunAtLoad</key>
    <false/>
</dict>
</plist>
